# supervisor.conf
## This is the main process for the Supervisor
[supervisord]
## This setting is to specify that we are not running in daemon mode
nodaemon = true

## This is the part where we give the name and add config for our 1st service
[program:citefix-bkend]
file=/tmp/supervisor.sock
chmod=0766
## This is the main command to run our 1st service
command = uvicorn main:app --reload --port 8000
## This setting specifies that the supervisor will restart the service in case of failure
autorestart = false
## This setting specifies that the supervisor will log the errors in the standard output
stderr_logfile = /dev/stdout
stderr_logfile_maxbytes = 0
## This setting specifies that the supervisor will log the output in the standard output
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
user = app